# wasde-monitor.service
[Unit]
Description=WASDE Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/wasde-monitor
ExecStart=/home/wasde-monitor/venv/bin/python /home/wasde-monitor/scheduler.py

Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment variables (alternatively, use EnvironmentFile=)
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
